# etc/systemd/system/influxdb.service.d/override.conf
[Service]
ExecStartPost=/bin/sh -c 'while ! curl --get --fail --silent http://localhost:8086/ping >/dev/null 2>&1; do sleep 1; done'
TimeoutStartSec=900
