# etc/systemd/system/influxdb.service.d/override.conf
[Service]
ExecStartPost=/bin/sh -c 'while ! /usr/bin/influx -execute exit >/dev/null 2>&1; do sleep 1; done'
TimeoutStartSec=900
