# etc/systemd/system/environment-monitor.target
[Unit]
Description=Raspberry Pi environment monitoring
Wants=envlog-monitor.service

[Install]
WantedBy=multi-user.target
